# ansible/files/xserver.service
[Unit]
Description=Standalone X

[Service]
Type=simple
ExecStart=/usr/bin/Xorg :0 -nolisten tcp -nocursor

[Install]
WantedBy=multi-user.target